# debian/python3-jobmanager-api.prerm
#!/usr/bin/env bash

/etc/init.d/jobmanager-api stop
update-rc.d -f jobmanager-api remove
